# exercise2/cleanup-logs.cron
#!/bin/bash
# Archivo: /etc/cron.d/cleanup-logs
# Descripción: Ejecuta el script cleanup_logs.sh cada día a las 2 AM
# Instalación: sudo cp cleanup-logs /etc/cron.d/

# Ejecutar cleanup_logs.sh cada día a las 2:00 AM
0 2 * * * root /home/devops/scripts/cleanup_logs.sh >> /var/log/cleanup_logs.log 2>&1
